# docker/nginx/crontab
52 3,15 * * * root /opt/certbot/renew >> /dev/null
