# docker/nginx/crontab
52 3,15 * * * /opt/certbot/renew >> /var/log/certbot/renew.log 
